# configs/etc/profile.d/julia.sh
#!/usr/bin/env sh

export JULIA_NUM_THREADS=$(nproc)

