# configs/etc/cron.weekly/portage-sync
#!/bin/sh

su - portage -s /bin/bash -c '(cd /usr/portage/local/layman; ./pull_all compress)'
